# git_bridge/rootfs/etc/services.d/git-bridge/run
#!/usr/bin/with-contenv sh
set -e

# Start FastAPI (bind to all interfaces inside the container)
# Use absolute path to avoid PATH/env surprises
exec /opt/venv/bin/uvicorn app:app --host 0.0.0.0 --port 8000
