# git_bridge/rootfs/etc/services.d/git-bridge/run
#!/bin/sh
set -eu

echo "=== git-bridge: starting ==="

# Ensure we run from /app (your app.py and requirements.txt are there)
cd /app

# Make sure Python can import from /app (so `import app` works)
export PYTHONPATH="/app:${PYTHONPATH:-}"

# Start uvicorn using the venv python (most reliable in HA add-ons)
exec /opt/venv/bin/python -m uvicorn app:app --host 0.0.0.0 --port 8000
