# git_bridge/rootfs/etc/services.d/git-bridge/run
#!/usr/bin/with-contenv bash
set -e

echo "=== git-bridge startup ==="
echo "whoami:"
whoami
echo "pwd:"
pwd

echo "ls /app:"
ls -la /app

echo "Setting PYTHONPATH=/app"
export PYTHONPATH=/app

echo "Testing python import:"
python3 - << 'EOF'
import sys
print("sys.path =", sys.path)
import app
print("Imported app OK")
print("FastAPI app object:", getattr(app, "app", None))
EOF

echo "Starting uvicorn (this must stay running)..."
cd /app

exec python3 -m uvicorn app:app \
  --host 0.0.0.0 \
  --port 8000
