# git_bridge/rootfs/etc/services.d/git-bridge/run
#!/bin/sh
set -e

# Always run uvicorn using the venv python to guarantee deps are available
exec /opt/venv/bin/python -m uvicorn app:app --host 0.0.0.0 --port 8000
