# git_bridge/rootfs/etc/services.d/git-bridge/run
#!/bin/sh
set -e

# Start the API (must stay in foreground; exec is required for s6 supervision)
exec python3 /app/app.py
