# git_bridge/rootfs/etc/services.d/git-bridge/run
# 1) Stop the add-on so it isn't holding files open
ha addons stop 95786995_git_bridge

# 2) Find where Home Assistant stored the add-on source on disk
#    (This is the repo path that contains rootfs/etc/services.d/.../run)
echo "=== Searching for the add-on folder that contains rootfs/etc/services.d ==="
find /data -maxdepth 6 -type f -path "*/rootfs/etc/services.d/*/run" 2>/dev/null | grep -E "95786995|git-bridge|git_bridge" || true
find / -maxdepth 8 -type f -path "*/rootfs/etc/services.d/*/run" 2>/dev/null | grep -E "95786995|git-bridge|git_bridge" || true

# 3) Print all candidates so you can see them clearly
echo "=== All candidate run files ==="
find /data -type f -path "*/rootfs/etc/services.d/*/run" 2>/dev/null | sort

# 4) Pick the one that matches your add-on and show its current content
#    (If multiple match, we'll overwrite the correct one next step)
echo "=== Showing current content of ALL candidate run files (first 80 lines) ==="
for f in $(find /data -type f -path "*/rootfs/etc/services.d/*/run" 2>/dev/null); do
  echo "----- $f -----"
  sed -n '1,80p' "$f" || true
done
