# git_bridge/rootfs/etc/services.d/git-bridge/run
#!/usr/bin/with-contenv sh
set -e

exec python3 /app/app.py
