# git_bridge/rootfs/etc/services.d/git-bridge/run
#!/bin/sh
set -e

echo "=== git-bridge boot ==="
echo "PWD=$(pwd)"
echo "--- /app listing ---"
ls -ლა /app || true

echo "--- venv python ---"
/opt/venv/bin/python -V || true

echo "--- pip freeze (first 50) ---"
/opt/venv/bin/pip freeze | head -n 50 || true

echo "--- test import ---"
/opt/venv/bin/python -c "import sys; sys.path.insert(0,'/app'); import app; print('IMPORTED OK', app.app)" || true

echo "--- starting uvicorn ---"
exec /opt/venv/bin/python -m uvicorn app:app --app-dir /app --host 0.0.0.0 --port 8000
