# git_bridge/rootfs/etc/services.d/git-bridge/run
#!/usr/bin/with-contenv sh
set -e

# Ensure the module path is correct for imports
cd /app
export PYTHONPATH="/app"

# Run uvicorn from the venv (don’t rely on PATH)
exec /opt/venv/bin/python -m uvicorn app:app --host 0.0.0.0 --port 8000
